# Week12/init/javaExample/javaservice.service
[Unit]
Description=A java service
After=hello.service

[Service]
ExecStart=/opt/startService.sh

[Install]
WantedBy=multi-user.target
